# matomo/systemd/matomo-archive@.timer
[Unit]
Description=Run Matomo archive every 5 minutes (%i)

[Timer]
OnBootSec=10min
OnActiveSec=10min
OnUnitActiveSec=30min
Unit=matomo-archive@%i.service

[Install]
WantedBy=timers.target
