# matomo/systemd/matomo-archive@.timer
[Unit]
Description=Run Matomo archive every 5 minutes (%i)

[Timer]
OnBootSec=5min
OnActiveSec=5min
Unit=matomo-archive@%i.service

[Install]
WantedBy=timers.target
